# deploy/jelleo-health.service
[Unit]
Description=Jelleo - daemon health check (one-shot)
Documentation=https://github.com/Copenhagen0x/audit-pipeline-cli
# Health-check crashes are invisible without this. Exit code 2 means
# degradation per `audit-pipeline health` — that's the EXPECTED unhappy
# path and is captured via webhook inside the command itself. OnFailure
# here covers the unexpected: the command itself crashing.
OnFailure=jelleo-alert-failure@%n.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/audit_runs/percolator-live
ExecStart=/bin/bash -c 'source /root/.audit-env && exec audit-pipeline --workspace /root/audit_runs/percolator-live health --silent'
StandardOutput=journal
StandardError=journal
