# deploy/jelleo-scheduler-monthly.service
[Unit]
Description=Jelleo - monthly cadence report dispatch
Documentation=https://jelleo.com/methodology.html#reporting
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Environment=PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=/root/audit_runs/percolator-live
ExecStart=/bin/bash -c 'source /root/.audit-env && exec audit-pipeline --workspace /root/audit_runs/percolator-live scheduler tick --cadence monthly'
StandardOutput=append:/root/audit_runs/percolator-live/scheduler/monthly.log
StandardError=append:/root/audit_runs/percolator-live/scheduler/monthly.log
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
