# deploy/jelleo-backup.service
[Unit]
Description=Jelleo - daily findings DB backup
Documentation=https://github.com/Copenhagen0x/audit-pipeline-cli
# Silent backup failures = data-retention SLA breach. Alert always.
OnFailure=jelleo-alert-failure@%n.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/audit_runs/percolator-live
ExecStart=/bin/bash /root/audit-pipeline-cli/deploy/backup_findings_db.sh
StandardOutput=journal
StandardError=journal
