# deploy/jelleo-scheduler-weekly.timer
[Unit]
Description=Jelleo - weekly cadence trigger
Documentation=https://jelleo.com/methodology.html#reporting

[Timer]
# Mondays at 09:15 UTC. 15-min offset from 24h to avoid simultaneous fire.
OnCalendar=Mon *-*-* 09:15:00 UTC
Persistent=true
RandomizedDelaySec=120
Unit=jelleo-scheduler-weekly.service

[Install]
WantedBy=timers.target
