# deploy/jelleo-watch.service
[Unit]
Description=Jelleo - source-code watch + autonomous hunt loop (Layer 0.5)
Documentation=https://github.com/Copenhagen0x/audit-pipeline-cli
After=network-online.target
Wants=network-online.target
# Cross-cutting audit Defect 24 (LOW operational): if the watch loop
# crashes hard enough to exhaust Restart=on-failure retries, the
# pipeline goes silent. OnFailure fires jelleo-alert-failure@.service
# which posts the unit name + journalctl tail to the alert channel
# the operator already has wired for health degradation.
OnFailure=jelleo-alert-failure@%n.service

[Service]
Type=simple
User=root
Environment=PATH=/root/.cargo/bin:/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=/root/audit_runs/percolator-live
ExecStartPre=/bin/mkdir -p /root/audit_runs/percolator-live/watch
# Watch polls upstream every 5 min. On new commit, runs a CHEAP triage cycle
# (recon-only, skip-poc / skip-debate / skip-narrative / skip-kani /
# skip-propagate / skip-bundle) to catch obvious regressions fast. Full
# Step-3-class audits are user-triggered explicitly with no skips.
#
# IMPORTANT: --on-update is wrapped in `bash -lc '...'` so shell builtins
# like `source` work. The previous form (`source ... && audit-pipeline ...`
# without bash wrap) failed silently with "[Errno 2] No such file or
# directory: 'source'" because watch.py exec()'s the command directly.
ExecStart=/bin/bash -c 'source /root/.audit-env && exec audit-pipeline --workspace /root/audit_runs/percolator-live watch --source-mode --interval 300 --on-update "/root/audit-pipeline-cli/deploy/watch_on_update.sh"'
Restart=on-failure
RestartSec=30s
StandardOutput=append:/root/audit_runs/percolator-live/watch/daemon.log
StandardError=append:/root/audit_runs/percolator-live/watch/daemon.log
KillMode=mixed
TimeoutStopSec=60s

[Install]
WantedBy=multi-user.target
