# deploy/jelleo-shadow.service
[Unit]
Description=Jelleo - live mainnet shadow audit (Layer 6)
Documentation=https://github.com/Copenhagen0x/audit-pipeline-cli
After=network-online.target
Wants=network-online.target
# Cross-cutting audit Defect 24 (LOW operational): alert on hard failure.
OnFailure=jelleo-alert-failure@%n.service

[Service]
Type=simple
User=root
Environment=PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=/root/audit_runs/percolator-live
ExecStartPre=/bin/mkdir -p /root/audit_runs/percolator-live/shadow
ExecStart=/bin/bash -c 'source /root/.audit-env && exec audit-pipeline --workspace /root/audit_runs/percolator-live shadow start --program 6qWZvUtfyShbxTQkwjCayk3LuGqTGJwBo2QfkePK5jdJ --watch-account CJKBStEn5VXEF9VNTChKKb5YW84MV7LycqMMziVuxJSc --interval 60 --limit 25'
Restart=on-failure
RestartSec=15s
StandardOutput=append:/root/audit_runs/percolator-live/shadow/daemon.log
StandardError=append:/root/audit_runs/percolator-live/shadow/daemon.log
KillMode=mixed
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
